# infra/supervisord.conf
[program:daemon-910985]
directory=/home/forge/bdi2.jerome-arfouche.ca/current/
command=bash -c 'set -a; source .env; set +a; exec node server.js'
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
user=forge
numprocs=1
redirect_stderr=true
stdout_logfile=/home/forge/.forge/daemon-910985.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=3
stopwaitsecs=15
stopasgroup=true
killasgroup=true
; Set these in Forge's environment panel — they are picked up automatically.
; If Forge does not inject them, define them here instead (no spaces around =):
; environment=DATABASE_URL="postgres://...",SESSION_SECRET="...",NODE_ENV="production"